# scripts/99-kfmon.rules
# Runs early at boot... (onboard *might* be mounted at that point)
# It'll run after on-animator, but regardless, because of our history, that's the perfect place to handle update cleanups.
KERNEL=="loop0", RUN+="/usr/local/kfmon/bin/kfmon-update.sh"
